# base-files/basefiles/profile.d/99-pipewire_alsa-fix.sh
# Fix alsa for programs like davinci resolve
export PIPEWIRE_ALSA='{ alsa.buffer-bytes=131072 }'
