# pika-sources/pika-sources/debian/postinst
#!/bin/sh

set -e

systemctl enable --now steam-repos-fix.service || echo "steam-repos-fix service could not be enabled"

