# auto-cpufreq/auto-cpufreq-1.9.6/debian/postinst
#!/bin/sh

set -e

systemctl enable --now auto-cpufreq.service || echo "auto-cpufreq service could not be enabled"

