# auto-cpufreq/auto-cpufreq-1.9.6/debian/prerm
#!/bin/sh

set -e

systemctl disable --now auto-cpufreq.service

