# wine-meta/wine-meta/debian/postinst
#!/bin/sh

set -e

systemctl enable --now link-opt-wine-to-bin.service || echo "link-opt-wine-to-bin service could not be enabled"

