# wine-meta/wine-meta/debian/prerm
#!/bin/sh

set -e

systemctl disable --now link-opt-wine-to-bin.service || echo "link-opt-wine-to-bin service could not be disabled"

