# pika-baseos/pika-baseos/etc/profile.d/99-pipewire_alsa-fix.sh
# Fix alsa for programs like davinci resolve, this fixes the Davinci Resolve audio delay
export PIPEWIRE_ALSA='{ alsa.buffer-bytes=20480 }'
