# nobara-session/nobara-session/etc/profile.d/nobara-gnome-libunity.sh
# Apply libunity env vars to nobara session
export XDG_CURRENT_DESKTOP="ubuntu:$XDG_CURRENT_DESKTOP" 
export XDG_SESSION_DESKTOP="ubuntu:$XDG_SESSION_DESKTOP" 
export DESKTOP_SESSION="ubuntu:$DESKTOP_SESSION"
