# wine-meta/wine-meta/etc/systemd/system/link-opt-wine-to-bin.service
[Unit]
Description=Links binaries from opt path to bin

[Service]
ExecStart=/usr/bin/wine-linker

[Install]
WantedBy=multi-user.target
