# pika-baseos/pika-baseos/etc/xdg/autostart/supergfxctl-gex-reset.desktop
[Desktop Entry]
Icon=device_pci
Name=Reset SuperGFXCTL-GEX
Exec=/etc/supergfxctl-gex-reset
Terminal=false
Type=Application
NoDisplay=true
